# config/supervisor_config/easyfinance.conf
[program:easyfinance]
command= /home/easyfinance-api/app/start_production.sh
directory= /home/easyfinance-api/app
user= easyfinance-api