# config/supervisor_config/easyfinance.conf
[program:storyexchange]
command= /home/deploy/easyfinance-api/scripts/start_production.sh
directory= /home/deploy/easyfinance-api/scripts
user= deploy